# deployment-files/etc/systemd/system/starburger-clearsessions.timer
[Unit]
Description=Timer for Starburger clearsessions

[Timer]
OnBootSec=300
OnUnitActiveSec=1w

[Install]
WantedBy=multi-user.target
